# deploy/systemd/user/splatlab.service.d/90-flight-a-recovery.conf
[Unit]
# A failed recovery leaves splatlab stopped. This prevents a marker-free service
# start after a power loss in the narrow archive-to-relock transition window.
Requires=splatlab-flight-a-boot-recovery.service
After=splatlab-flight-a-boot-recovery.service

